# ansible-artifactory/roles/install-artifactory/files/artifactory.ini
[program:artifactory]
command=service artifactory start
user=root
stdout_logfile=/var/log/supervisor/artifactory.log
stderr_logfile=/var/log/supervisor/artifactory.err
